# boot_dump/ramdisk/system/etc/ld.so.preload
/lib/libc_stubs.so
/lib/libc++.so.1.0
/lib/libc++abi.so.1.0
